# Socat-Expand/files/Socat-Expand
#!/bin/sh /etc/rc.common

START=99

start()
{
	[ -f /etc/socat ] && mv /etc/socat /etc/config/socat
	rm -f /etc/config/socat-okpg
}
